# configs/com.vitaclaw.heartbeat.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vitaclaw.heartbeat</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/env</string>
        <string>python3</string>
        <string>scripts/run_health_heartbeat.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>__VITACLAW_ROOT__</string>
    <key>StartInterval</key>
    <integer>7200</integer>
    <key>StandardOutPath</key>
    <string>__VITACLAW_ROOT__/logs/heartbeat.log</string>
    <key>StandardErrorPath</key>
    <string>__VITACLAW_ROOT__/logs/heartbeat.err</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
